# config/com.clawy.bettingbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawy.bettingbot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/olli/.openclaw/workspace/multi-sport-signal-bot/.venv/bin/python</string>
        <string>-c</string>
        <string>import sys; sys.path.insert(0, '/Users/olli/.openclaw/workspace/multi-sport-signal-bot'); from src.bot.app import main; main()</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/olli/.openclaw/workspace/multi-sport-signal-bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/olli/.openclaw/workspace/multi-sport-signal-bot/logs-bot.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/olli/.openclaw/workspace/multi-sport-signal-bot/logs-bot.err</string>
</dict>
</plist>
